# apps/beralyzer/beralyzer.service.example
[Unit]
Description=Beralyzer daemon (mainnet)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
Environment=PG_DSN=postgres://beralyzer:***@localhost:5432/beralyzer_mainnet
Environment=EL_ETHRPC_URL=http://127.0.0.1:8545
Environment=CL_ETHRPC_URL=http://127.0.0.1:26657
Environment=BERALYZER_POLL_MS=15000
WorkingDirectory=/storage/guides/exp/beralyzer
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
User=beralyzer

[Install]
WantedBy=multi-user.target


